# systemd/user/voxtype.service.d/gpu.conf
[Service]
Environment="VOXTYPE_VULKAN_DEVICE=amd"
