# autostart/autostart/Handy.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Handy
Comment=Handystartup script
Exec=/usr/bin/handy 
StartupNotify=false
Terminal=false